# execd/figlet/execd.service
[Unit]
Description=execd

[Service]
Type=simple
User=nobody
ExecStart=/usr/local/bin/execd -listen :4000

[Install]
WantedBy=multi-user.target
